# templates/ninjablock.erb
#!/bin/sh -e
#
# /etc/init.d/ninjablock: ninjablock client script for Debian
#
### BEGIN INIT INFO
# Provides:   ninjablock
# Required-Start: $local_fs $syslog $remote_fs
# Required-Stop: $remote_fs
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start ninjablock client
# Description: starts ninjablock client
### END INIT INFO

NAME=ninjablock
PATH_BIN=/bin:/usr/bin:/sbin:/usr/sbin
DAEMON=/usr/bin/nodejs
LOG_FILE=/var/log/ninjablock.log
OPTS=<%= @dir %>/client.js
PIDFILE=/var/run/$NAME.pid
DEFAULTS=/etc/default/$NAME
DESC="ninjablock daemon"
USER=<%= @user %>

ENV="env -i LANG=C PATH=$PATH_BIN"

. /lib/lsb/init-functions

test -x $DAEMON || exit 0

RET=0

[ -r "$DEFAULTS" ] && . "$DEFAULTS" || start_ninjablock_daemon=yes


check_daemon_enabled () {
    if [ "$start_ninjablock_daemon" = "yes" ] ; then
        return 0
    else
        log_warning_msg "Not starting $NAME daemon, disabled via /etc/default/$NAME"
        return 1
    fi

}


case "$1" in
  start)
    if check_daemon_enabled ; then
      log_daemon_msg "Starting $DESC" "$NAME" 
      touch $LOG_FILE && chown $USER:$USER $LOG_FILE
      if start-stop-daemon --start --oknodo --chuid $USER -b --make-pidfile --pidfile $PIDFILE --exec $DAEMON -- $OPTS; then
        log_end_msg 0
      else
        log_end_msg 1
        RET=1
      fi
    fi
    ;;

  stop)
    log_daemon_msg "Stopping $DESC" "$NAME"
    if start-stop-daemon --stop --oknodo --pidfile $PIDFILE --retry 30 ; then
        log_end_msg 0
    else
        log_end_msg 1
        RET=1
    fi
    ;;

  restart)
    if check_daemon_enabled ; then
        log_action_begin_msg "Restarting $NAME"
        if $0 stop && $0 start ; then
            log_action_end_msg 0
        else
            log_action_cont_msg "(failed)"
            RET=1
        fi
    fi
    ;;

  status)
    start-stop-daemon --status --pidfile $PIDFILE
    ;;

  *)
    echo "Usage: /etc/init.d/$NAME {start|stop|status|restart}"
    RET=1
    ;;
esac


exit $RET
